# agents/systemd/montime-agent-bash.service
[Unit]
Description=Montime.io Monitoring Agent (Bash)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/montime
Environment="SERVER_TOKEN=YOUR_SERVER_TOKEN_HERE"
Environment="BASE_URL=https://montime.io"
ExecStart=/bin/bash /opt/montime/agent.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=montime-agent

[Install]
WantedBy=multi-user.target
